# roles/rc/files/rc.local
#!/bin/sh -e
#
# rc.local -- Added for compatibility with systemd
#
# Enable with:
#
# $ systemctl unmask rc-local
# $ systemctl enable rc-local
# $ systemctl start rc-local
#

exit 0

